# utils/ath9k-watchdog/files/ath9k-watchdog.init
#!/bin/sh /etc/rc.common

START=95
APP=ath9k-watchdog
PID_FILE=/var/run/$APP.pid

start() {
		if [ -f $PID_FILE ]; then
			echo "Already Running"
			exit 1
		fi
        start-stop-daemon -S -x $APP -p $PID_FILE -m -b
}

stop() {
        start-stop-daemon -K -n $APP -p $PID_FILE -s TERM
        rm -rf $PID_FILE
}